# overlay.d/20platform-chrony/usr/lib/systemd/system/chronyd.service.d/platform-chrony.conf
# Override chrony configuration with the output of
# coreos-platform-chrony-config.service.
[Service]
EnvironmentFile=-/run/nestos/sysconfig-chrony
